# files/etc/uci-defaults/99-ath10k-aql-disable
#!/bin/sh

/etc/init.d/ath10k-aql-disable enable
/etc/init.d/ath10k-aql-disable start

exit 0
