# files/etc/init.d/ath10k-aql-disable
#!/bin/sh /etc/rc.common

START=99

# Test mitigation for ath10k "failed to lookup txq" stalls/disconnects:
# https://github.com/openwrt/openwrt/issues/9455
# The debugfs knobs can appear after init starts, so retry briefly.
apply_aql_mitigation()
{
    local aql found

    found=0
    for aql in /sys/kernel/debug/ieee80211/phy*/aql_enable
    do
        [ -e "$aql" ] || continue
        echo 0 > "$aql"
        found=1
        logger -t ath10k-aql-disable "disabled AQL via $aql"
    done

    [ "$found" -eq 1 ]
}

start()
{
    local retry

    for retry in 1 2 3 4 5 6 7 8 9 10
    do
        apply_aql_mitigation && return 0
        sleep 3
    done

    logger -t ath10k-aql-disable "no AQL debugfs controls found"
    return 0
}
